# ansible/files/borg-backup.timer
[Unit]
Description=Borg Backup

[Timer]
OnUnitActiveSec=1min
Unit=borg-backup.service

[Install]
WantedBy=timers.target
